# sentinel-rules/scheduled-task-persistence.yaml
id: native-scheduled-task-persistence-analytic-rule
kind: Scheduled
name: NATIVE – Scheduled Task Persistence
description: |
  L3-native analytic rule for Scheduled Task Persistence, built from the Native Detection Bible.
  Uses only native telemetry and includes ConfidenceScore, Severity and HuntingDirectives.

severity: High
requiredDataConnectors:
  - connectorId: MicrosoftThreatProtection
    dataTypes:
      - DeviceProcessEvents

query: |
  // See rules/persistence/scheduled-task-persistence.md for full KQL with scoring and directives.

queryFrequency: 1h
queryPeriod: 1h
triggerOperator: GreaterThan
triggerThreshold: 0

suppressionEnabled: false
suppressionDuration: 6h

tactics:
  - execution

techniques:
  - T1053.002

entityMappings:
  - entityType: Host
    fieldMappings:
      - identifier: FullName
        columnName: DeviceName
  - entityType: Account
    fieldMappings:
      - identifier: Name
        columnName: AccountName
